# roles/music_server_snapcast/files/bbc.service
[Unit]
Description=Created BBC Radio Playlists

[Service]
Type=simple
ExecStart=/opt/bbc.sh
