# roles/doorbell/templates/doorbell.service
[Unit]
Description=doorbell
After=network.target

[Service]
User=root
ExecStart=/usr/bin/python /opt/doorbell/pi/doorbell.py

[Install]
WantedBy=default.target
